# empire/server/modules/powershell/persistence/elevated/schtasks.yaml
name: Invoke-Schtasks
authors:
  - name: Matt Graeber
    handle: '@mattifestation'
    link: https://twitter.com/mattifestation
  - name: Will Schroeder
    handle: '@harmj0y'
    link: https://twitter.com/harmj0y
description: |
  Establishes persistence by creating a scheduled task that runs with
  SYSTEM privileges to execute a stager or script. This module leverages the Windows
  Task Scheduler to create a reliable persistence mechanism that can trigger on
  various events including daily execution, user logon, or idle time.

  Scheduled task persistence provides moderate stealth as it uses legitimate Windows functionality
  and can be configured with various triggers to blend with normal system activity.
  The module can store the payload in the registry or reference an external file,
  and supports multiple scheduling options for operational flexibility.

  This
  technique requires administrative privileges and creates visible scheduled task
  artifacts that can be detected by security monitoring tools.
software: S0111
tactics: [TA0003, TA0004]
techniques: [T1053.005]
background: false
output_extension:
needs_admin: true
opsec_safe: false
language: powershell
min_language_version: '2'
comments:
  - https://github.com/mattifestation/PowerSploit/blob/master/Persistence/Persistence.psm1
options:
  - name: Listener
    description: Listener to use for the stager.
    required: false
    value: ''
  - name: Obfuscate
    description: Obfuscate the launcher PowerShell code to evade signature-based
      detection.
    required: false
    value: 'False'
    strict: true
    suggested_values:
      - True
      - False
  - name: ObfuscateCommand
    description: The Invoke-Obfuscation command to use.
    required: false
    value: Token\All\1
    depends_on:
      - name: Obfuscate
        values: ['True']
  - name: Bypasses
    description: Bypasses as a space separated list to be prepended to the launcher.
    required: false
    value: ''
  - name: DailyTime
    description: 'Daily time to trigger the script (format: HH:mm).'
    required: false
    value: 09:00
  - name: IdleTime
    description: User idle time in minutes to trigger script (leave blank to disable).
    required: false
    value: ''
  - name: OnLogon
    description: Trigger script on user logon (set to True to enable).
    required: false
    value: 'False'
    strict: true
    suggested_values:
      - True
      - False
  - name: TaskName
    description: Name for the scheduled task.
    required: true
    value: Updater
  - name: RegPath
    description: Registry location to store the script code (last element is the
      key name).
    required: false
    value: HKLM:\Software\Microsoft\Network\debug
  - name: ADSPath
    description: Alternate-data-stream location to store the script code (leave blank
      to use registry).
    required: false
    value: ''
  - name: ExtFile
    description: Path to external file containing the payload (leave blank to use
      stager).
    required: false
    value: ''
  - name: Cleanup
    description: Remove the persistence mechanism and clean up artifacts.
    required: false
    value: 'False'
    strict: true
    suggested_values:
      - True
      - False
  - name: UserAgent
    description: User-agent string to use for the staging request (default, none,
      or other).
    required: false
    value: default
  - name: Proxy
    description: Proxy to use for request (default, none, or other).
    required: false
    value: default
  - name: ProxyCreds
    description: Proxy credentials ([domain\]username:password) to use for request
      (default, none, or other).
    required: false
    value: default
advanced:
  custom_generate: true
