# empire/server/modules/powershell/persistence/userland/registry.yaml
name: Invoke-Registry
authors:
  - name: Matt Graeber
    handle: '@mattifestation'
    link: https://twitter.com/mattifestation
  - name: Will Schroeder
    handle: '@harmj0y'
    link: https://twitter.com/harmj0y
  - name: ''
    handle: '@enigma0x3'
    link: ''
description: |
  Establishes persistence by creating a registry entry in HKCU:SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run
  that executes a stager or script during user logon. This module creates a user-level
  persistence mechanism that triggers when the specific user logs into the system.

  User-level registry persistence is less privileged than system-level persistence but provides
  reliable execution for the target user account. The module can store the payload
  in the registry, alternate data streams, event logs, or reference an external
  file, providing flexibility for different operational requirements and stealth
  considerations.

  This technique does not require administrative privileges and
  creates visible registry artifacts that can be detected by security monitoring
  tools.
software: ''
tactics: [TA0003]
techniques: [T1547.001]
background: false
output_extension:
needs_admin: false
opsec_safe: false
language: powershell
min_language_version: '2'
comments:
  - https://github.com/mattifestation/PowerSploit/blob/master/Persistence/Persistence.psm1
options:
  - name: Listener
    description: Listener to use for the stager.
    required: false
    value: ''
  - name: Obfuscate
    description: Obfuscate the launcher PowerShell code to evade signature-based
      detection.
    required: false
    value: 'False'
    strict: true
    suggested_values:
      - True
      - False
  - name: ObfuscateCommand
    description: The Invoke-Obfuscation command to use.
    required: false
    value: Token\All\1
    depends_on:
      - name: Obfuscate
        values: ['True']
  - name: Bypasses
    description: Bypasses as a space separated list to be prepended to the launcher.
    required: false
    value: ''
  - name: KeyName
    description: Name for the registry key that will trigger execution.
    required: true
    value: Updater
  - name: RegPath
    description: Registry location to store the script code (last element is the
      key name).
    required: false
    value: HKCU:Software\Microsoft\Windows\CurrentVersion\Debug
  - name: ADSPath
    description: Alternate-data-stream location to store the script code (leave blank
      to use registry).
    required: false
    value: ''
  - name: EventLogID
    description: Store the script in the Application event log under the specified
      EventID (must be unique/rare).
    required: false
    value: ''
  - name: ExtFile
    description: Path to external file containing the payload (leave blank to use
      stager).
    required: false
    value: ''
  - name: Cleanup
    description: Remove the persistence mechanism and clean up artifacts.
    required: false
    value: 'False'
    strict: true
    suggested_values:
      - True
      - False
  - name: UserAgent
    description: User-agent string to use for the staging request (default, none,
      or other).
    required: false
    value: default
  - name: Proxy
    description: Proxy to use for request (default, none, or other).
    required: false
    value: default
  - name: ProxyCreds
    description: Proxy credentials ([domain\]username:password) to use for request
      (default, none, or other).
    required: false
    value: default
advanced:
  custom_generate: true
